# stow/macos/Library/LaunchAgents/com.vili.homerow.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.vili.homerow</string>
	<!-- Homerow ignores open's -j (launch hidden) hint and shows its settings
	     window on every launch. So: launch it, wait for the process, then hide
	     the app via System Events (equivalent to Cmd+H — doesn't need
	     accessibility access). Hide twice in case the window appears late.
	     NOTE: first login after installing this, macOS shows a one-time
	     "wants to control System Events" prompt — click Allow. -->
	<key>ProgramArguments</key>
	<array>
		<string>/bin/sh</string>
		<string>-c</string>
		<string>open -g -j -a Homerow
i=0; until pgrep -xq Homerow || [ $i -ge 20 ]; do sleep 0.5; i=$((i+1)); done
sleep 1
osascript -e 'tell application "System Events" to set visible of process "Homerow" to false' || true
sleep 2
osascript -e 'tell application "System Events" to set visible of process "Homerow" to false' || true</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
